# roles/kafka/templates/kafka.service.j2
[Unit]
Description=An Apache Kafka broker
After=network.target

[Service]
Type=forking
ExecStart=/usr/bin/kafka-server-start -daemon /etc/kafka/server.properties
ExecStop=/usr/bin/kafka-server-stop

